# recipes-core/udev/udev-rules-wpi/11-wpi.rules
SUBSYSTEM=="extcon", ACTION=="change", RUN+="/etc/udev/scripts/wpi_usb_en.sh"
